# make-openwrt/openwrt-files/common-files/etc/profile.d/sys_bashrc.sh
[ -z "$BASH" ] || [ "$BASH" = /bin/bash ] || [ \! -s /etc/bash.bashrc ] || . /etc/bash.bashrc
