# make-openwrt/openwrt-files/different-files/king3399/rootfs/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

bash /etc/custom_service/start_service.sh
sleep 1

bash /etc/init.d/king-gpio enable

exit 0
